# crontab
0 * * * * cd /post_processor && /usr/local/bin/python main.py --env prod >> /var/log/cron.log 2>&1
